# mt5_container/supervisor.conf
[supervisord]
nodaemon=true
logfile=/home/trader/supervisord.log
pidfile=/home/trader/supervisord.pid

[program:xvfb]
command=Xvfb :1 -screen 0 1280x1024x24 -fbdir /var/tmp
autorestart=true
priority=100

[program:fluxbox]
command=fluxbox -display :1
autorestart=true
priority=200

[program:x11vnc]
command=x11vnc -display :1 -nopw -forever -shared -listen 0.0.0.0
autorestart=true
priority=300

[program:mt5]
command=wine "C:\Program Files\MetaTrader 5\terminal64.exe" /portable
autorestart=true
priority=400
stdout_logfile=/home/trader/mt5_stdout.log
stderr_logfile=/home/trader/mt5_stderr.log

[program:grpc_bridge]
command=wine C:\Python310\python.exe /home/trader/app/bridge.py
autorestart=true
priority=500
stdout_logfile=/home/trader/bridge_stdout.log
stderr_logfile=/home/trader/bridge_stderr.log
